# kokoro/src/service/x64/com.github.tibssy.kokoro-tts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.kokoro-tts-user</string>

    <key>Program</key>
    <string>${HOME}/bin/kokoro-tts/kokoro-tts-macos-x64.bin</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/com.example.kokoro-tts.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/com.example.kokoro-tts.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CONFIG_PATH</key>
        <string>${HOME}/Library/Application Support/kokoro-tts/config.toml</string>
        <key>MODEL_PATH</key>
        <string>${HOME}/Library/Application Support/kokoro-tts/models/kokoro-v1.0.onnx</string>
        <key>VOICES_PATH</key>
        <string>${HOME}/Library/Application Support/kokoro-tts/models/voices-v1.0.bin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>${HOME}</string>
</dict>
</plist>