# kokoro/src/service/x64/kokoro-tts.service
[Unit]
Description=Kokoro TTS Service (User)
After=default.target

[Service]
ExecStart=%h/.local/bin/kokoro-tts-linux-x64.bin
Restart=always
RestartSec=1
Environment=CONFIG_PATH=%h/.config/kokoro-tts/config.toml
Environment=MODEL_PATH=%h/.local/share/kokoro-tts/models/kokoro-v1.0.onnx
Environment=VOICES_PATH=%h/.local/share/kokoro-tts/models/voices-v1.0.bin

[Install]
WantedBy=default.target